# overlay/system/etc/profile.d/devicehacks.sh
#!/bin/sh

export QMLSCENE_DEVICE=haliumqsgcontext
